# deploy/karyra-imgproxy.service.example
[Unit]
Description=Karyra Spark imgproxy media optimizer container
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/karyra/imgproxy.env
ExecStartPre=-/usr/bin/podman rm -f karyra-imgproxy
ExecStart=/usr/bin/podman run --name karyra-imgproxy --rm --env-file /etc/karyra/imgproxy.env -p 127.0.0.1:8088:8080 ghcr.io/imgproxy/imgproxy:latest-amd64
ExecStop=/usr/bin/podman stop -t 10 karyra-imgproxy
Restart=always
RestartSec=3
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
